# snotify.service
[Unit]
Description=Snotify user daemon

[Service]
ExecStart=/usr/bin/snotify
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
